# supervisord.d/integral_rpc.ini
[program:integral_rpc]
command=/root/.jenkins/workspace/integral.rpc.server/integral/command/rpc/integral.rpc -f /root/.jenkins/workspace/integral.rpc.server/integral/command/rpc/config/config.json
directory=/root/.jenkins/workspace/integral.rpc.server/integral/command/rpc/
