# supervisord.d/user_rpc.ini
[program:user_rpc]
command=/root/.jenkins/workspace/user.api.server/user/command/rpc/user_rpc -f /root/.jenkins/workspace/user.api.server/user/command/rpc/config/config.json
directory=/root/.jenkins/workspace/user.api.server/user/command/rpc/
