# supervisord.d/user_api.ini
[program:user_api]
command=/root/.jenkins/workspace/user.api.server/user/command/api/user.api -f /root/.jenkins/workspace/user.api.server/user/command/api/config/config.json
directory=/root/.jenkins/workspace/user.api.server/user/command/api/
